# services/data_drift.service
[Unit]
Description=model_deployer
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=3
Environment="ENV=DEV"
WorkingDirectory=/home/vboxuser/BT4301_Project/services
ExecStart=/home/vboxuser/BT4301_Project/services/venv/bin/python3 /home/vboxuser/BT4301_Project/services/data_drift.py

[Install]
WantedBy=multi-user.target